# docker/scripts/crontab
10 0 * * 0 symfony /usr/bin/cron-script

